# sahara/plugins/vanilla/hadoop2/resources/spark-cleanup.cron
# Cleans up old Spark job directories once per hour.
0 * * * * root /etc/hadoop/tmp-cleanup.sh